# cron/cron.d/crontab
* * * * * /bin/bash /cron/run.sh > /var/log/finder

# Must have Blank line before EOF
